# ddp.service
[Unit]
Description=ddp response handler

[Service]
Type=simple
ExecStart=/usr/bin/ddpd

[Install]
WantedBy=multi-user.target
